# files/root/etc/uci-defaults/luci-ssr
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ssr[-1]
	add ucitrack ssr
	set ucitrack.@ssr[-1].init=ssr
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
